# scripts/services/rag_conversation.service
[Unit]
Description=T0 RAG Conversational Service
After=network.target t0_main.service qwen_with_tools.service
# Disable for now, as this seems too demanding
# Requires=t0_main.service qwen_with_tools.service

[Service]
Type=simple
User=t0
WorkingDirectory=/home/t0/t0-1
ExecStart=/home/t0/.local/bin/uv run --directory /home/t0/t0-1 /home/t0/t0-1/scripts/serve_rag_conversational.sh
# ExecStart=/bin/bash -c 'ping 8.8.8.8'
# wait for the string "Application startup complete." in logs
ExecStartPost=/bin/bash -c 'timeout 900 /bin/bash /home/t0/t0-1/scripts/is_rag_conversation_ready.sh'
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
